# watcher/guildpost-watcher.service
[Unit]
Description=GuildPost Minecraft Server Watcher
After=network.target

[Service]
Type=simple
User=watcher
Group=watcher
WorkingDirectory=/opt/watcher
ExecStart=/opt/watcher/watcher
Environment="SUPABASE_URL=https://wpxutsdbiampnxfgkjwq.supabase.co"
Environment="SUPABASE_SERVICE_KEY=YOUR_SERVICE_KEY_HERE"
Environment="PING_API_URL=https://guildpost.pages.dev"
Environment="POLL_INTERVAL=300"
Environment="WORKER_COUNT=5"
Environment="LISTEN_ADDR=:8080"
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=watcher

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/watcher

[Install]
WantedBy=multi-user.target
